# superset/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- from "superset/map.jinja" import server with context %}

{%- if service.name == "web" %}

[program:{{ service.type }}_webserver]
directory=/srv/superset
command=/srv/superset/bin/start-service.sh "superset runserver"
stdout_logfile=/var/log/superset/web.log
stderr_logfile=/var/log/superset/web_error.log
user=superset
autostart=true
autorestart=true

{% endif %}

{% if service.name == "worker" %}

[program:{{ service.type }}_worker]
directory=/srv/superset
command=/srv/superset/bin/start-service.sh "superset worker --workers {{ server.get("workers", 5) }}"
stdout_logfile=/var/log/superset/worker.log
stderr_logfile=/var/log/superset/worker_error.log
user=superset
autostart=true
autorestart=true

{% endif %}
